# experimental/ros/noetic/pre-patch/etc/catkin/profile.d/999.sdformat.bat
set SDF_PATH=c:\opt\ros\noetic\x64\share\sdformat\1.6
