# experimental/ros/melodic/patch/etc/catkin/profile.d/999.graphviz.bat
@echo off

set PATH=%PATH:C:\opt\ros\melodic\x64\tools\graphviz;=%
set "PATH=C:\opt\ros\melodic\x64\tools\graphviz;%PATH%"
