# ros/melodic/patch/etc/catkin/profile.d/999.sdformat.bat
set SDF_PATH=c:\opt\ros\melodic\x64\share\sdformat\1.6
